# perplexitytrader-webhook.service
[Unit]
Description=PerplexityTrader Webhook Server
After=network.target

[Service]
Type=simple
User=angle
WorkingDirectory=/home/angle/perpleixtytrader
ExecStart=/home/angle/perpleixtytrader/bluefin_env/bin/python /home/angle/perpleixtytrader/webhook_server.py
Restart=always
RestartSec=10
Environment="FLASK_ENV=production"
Environment="FLASK_DEBUG=false"
Environment="WEBHOOK_PORT=8000"
Environment="USE_NGROK=true"
Environment="NGROK_AUTHTOKEN=2tvsk5sHN6PQoJNvsKsSf1DUDMp_dVpHp5FFy1VpXxf4oyn8"
Environment="NGROK_DOMAIN=https://awake-drake-bursting.ngrok-free.app"

[Install]
WantedBy=multi-user.target 